# pokerbot.service
[Unit]
Description=Poker Staking Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=marcus
WorkingDirectory=/home/marcus/poker-bot
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
